# etc/systemd/system/linuxmuster-client-auth-ldap-autoconfiguration.service
[Unit]
Description=linuxmuster-cliet-auth LDAP autoconfiguration
Requires=networking.service
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/linuxmuster-client-auth-ldap-autoconfiguration


[Install]
WantedBy=multi-user.target
